# core/dbus/.footprint
drwxr-xr-x etc/
drwxr-xr-x etc/dbus-1/
-rw-r--r-- etc/dbus-1/session.conf
-rw-r--r-- etc/dbus-1/system.conf
drwxr-xr-x etc/rc.d/
-rwxr-xr-x etc/rc.d/dbus
lrwxrwxrwx lib
drwxr-xr-x lib64/
drwxr-xr-x usr/
drwxr-xr-x usr/bin/
-rwxr-xr-x usr/bin/dbus-cleanup-sockets
-rwxr-xr-x usr/bin/dbus-daemon
-rwxr-xr-x usr/bin/dbus-launch
-rwxr-xr-x usr/bin/dbus-monitor
-rwxr-xr-x usr/bin/dbus-run-session
-rwxr-xr-x usr/bin/dbus-send
-rwxr-xr-x usr/bin/dbus-test-tool
-rwxr-xr-x usr/bin/dbus-update-activation-environment
-rwxr-xr-x usr/bin/dbus-uuidgen
drwxr-xr-x usr/include/
drwxr-xr-x usr/include/dbus-1.0/
drwxr-xr-x usr/include/dbus-1.0/dbus/
-rw-r--r-- usr/include/dbus-1.0/dbus/dbus-address.h
-rw-r--r-- usr/include/dbus-1.0/dbus/dbus-bus.h
-rw-r--r-- usr/include/dbus-1.0/dbus/dbus-connection.h
-rw-r--r-- usr/include/dbus-1.0/dbus/dbus-errors.h
-rw-r--r-- usr/include/dbus-1.0/dbus/dbus-macros.h
-rw-r--r-- usr/include/dbus-1.0/dbus/dbus-memory.h
-rw-r--r-- usr/include/dbus-1.0/dbus/dbus-message.h
-rw-r--r-- usr/include/dbus-1.0/dbus/dbus-misc.h
-rw-r--r-- usr/include/dbus-1.0/dbus/dbus-pending-call.h
-rw-r--r-- usr/include/dbus-1.0/dbus/dbus-protocol.h
-rw-r--r-- usr/include/dbus-1.0/dbus/dbus-server.h
-rw-r--r-- usr/include/dbus-1.0/dbus/dbus-shared.h
-rw-r--r-- usr/include/dbus-1.0/dbus/dbus-signature.h
-rw-r--r-- usr/include/dbus-1.0/dbus/dbus-syntax.h
-rw-r--r-- usr/include/dbus-1.0/dbus/dbus-threads.h
-rw-r--r-- usr/include/dbus-1.0/dbus/dbus-types.h
-rw-r--r-- usr/include/dbus-1.0/dbus/dbus.h
lrwxrwxrwx usr/lib
drwxr-xr-x usr/lib64/
drwxr-xr-x usr/lib64/dbus-1.0/
drwxr-xr-x usr/lib64/dbus-1.0/include/
drwxr-xr-x usr/lib64/dbus-1.0/include/dbus/
-rw-r--r-- usr/lib64/dbus-1.0/include/dbus/dbus-arch-deps.h
lrwxrwxrwx usr/lib64/libdbus-1.so
lrwxrwxrwx usr/lib64/libdbus-1.so.3
-rwxr-xr-x usr/lib64/libdbus-1.so.3.14.3
drwxr-xr-x usr/lib64/pkgconfig/
-rw-r--r-- usr/lib64/pkgconfig/dbus-1.pc
drwxr-xr-x usr/libexec/
-rwsr-x--- usr/libexec/dbus-daemon-launch-helper
drwxr-xr-x usr/share/
drwxr-xr-x usr/share/dbus-1/
drwxr-xr-x usr/share/dbus-1/services/
-rw-r--r-- usr/share/dbus-1/session.conf
drwxr-xr-x usr/share/dbus-1/session.d/
drwxr-xr-x usr/share/dbus-1/system-services/
-rw-r--r-- usr/share/dbus-1/system.conf
drwxr-xr-x usr/share/dbus-1/system.d/
drwxr-xr-x var/
drwxr-xr-x var/lib/
drwxr-xr-x var/lib/dbus/
drwxr-xr-x var/run/
drwxr-xr-x var/run/dbus/
