# core/autofs/.footprint
drwxr-xr-x etc/
-rw-r--r-- etc/auto.master
-rw-r--r-- etc/auto.misc
-rwxr-xr-x etc/auto.net
-rwxr-xr-x etc/auto.smb
-rw-r--r-- etc/autofs.conf
-rw------- etc/autofs_ldap_auth.conf
drwxr-xr-x etc/default/
-rw-r--r-- etc/default/autofs
drwxr-xr-x etc/rc.d/
-rwxr-xr-x etc/rc.d/autofs
lrwxrwxrwx lib
drwxr-xr-x lib64/
drwxr-xr-x usr/
lrwxrwxrwx usr/lib
drwxr-xr-x usr/lib64/
drwxr-xr-x usr/lib64/autofs/
-rwxr-xr-x usr/lib64/autofs/lookup_dir.so
-rwxr-xr-x usr/lib64/autofs/lookup_file.so
lrwxrwxrwx usr/lib64/autofs/lookup_files.so
-rwxr-xr-x usr/lib64/autofs/lookup_hosts.so
-rwxr-xr-x usr/lib64/autofs/lookup_multi.so
lrwxrwxrwx usr/lib64/autofs/lookup_nis.so
-rwxr-xr-x usr/lib64/autofs/lookup_nisplus.so
-rwxr-xr-x usr/lib64/autofs/lookup_program.so
-rwxr-xr-x usr/lib64/autofs/lookup_userhome.so
-rwxr-xr-x usr/lib64/autofs/lookup_yp.so
-rwxr-xr-x usr/lib64/autofs/mount_afs.so
-rwxr-xr-x usr/lib64/autofs/mount_autofs.so
-rwxr-xr-x usr/lib64/autofs/mount_bind.so
-rwxr-xr-x usr/lib64/autofs/mount_changer.so
-rwxr-xr-x usr/lib64/autofs/mount_ext2.so
lrwxrwxrwx usr/lib64/autofs/mount_ext3.so
lrwxrwxrwx usr/lib64/autofs/mount_ext4.so
-rwxr-xr-x usr/lib64/autofs/mount_generic.so
-rwxr-xr-x usr/lib64/autofs/mount_nfs.so
lrwxrwxrwx usr/lib64/autofs/mount_nfs4.so
-rwxr-xr-x usr/lib64/autofs/parse_amd.so
-rwxr-xr-x usr/lib64/autofs/parse_sun.so
drwxr-xr-x usr/sbin/
-rwxr-xr-x usr/sbin/automount
drwxr-xr-x usr/share/
drwxr-xr-x usr/share/man/
drwxr-xr-x usr/share/man/man5/
-rw-r--r-- usr/share/man/man5/auto.master.5
-rw-r--r-- usr/share/man/man5/autofs.5
-rw-r--r-- usr/share/man/man5/autofs.conf.5
-rw-r--r-- usr/share/man/man5/autofs_ldap_auth.conf.5
drwxr-xr-x usr/share/man/man8/
-rw-r--r-- usr/share/man/man8/autofs.8
-rw-r--r-- usr/share/man/man8/automount.8
drwxr-xr-x var/
drwxr-xr-x var/run/
